# openwrt_defaults/package/base-files/files/etc/uci-defaults/14_tsmodem
#!/bin/sh

uci -q batch <<-EOF >/dev/null
    set system.@system[0].hostname='BITCORD-RTR-02'
    set system.@system[0].pubname='Роутер BITCORD RTR-02'
    set system.@system[0].timezone='MSK-3'
    set system.@system[0].zonename='Europe/Moscow'
    commit system
    set wireless.default_radio0.key='WFCLokm7Ls'
    set wireless.default_radio0.ssid='BITCORD-RTR-02'
    set wireless.default_radio0.encryption='psk2'
    set wireless.radio0.disabled='0'
    commit wireless
    set dhcp.lan.ra_management='1'
    commit dhcp
EOF

uci -q batch <<-EOF
	set network.tsmodem=interface
	set network.tsmodem.username='root'
	set network.tsmodem.proto='3g'
	set network.tsmodem.device='/dev/ttyUSB2'
	set network.tsmodem.service='umts'
	set network.tsmodem.ipv6='auto'
	set network.tsmodem.apn='internet.mts.ru'
	delete network.lan.ifname='eth0.1'
	set network.lan.ifname='3g-tsmodem eth0.1'
	commit network
EOF

uci -q batch <<-EOF
	add_list firewall.@zone[1].network='tsmodem'
	delete firewall.@defaults[0].syn_flood='1'
	delete firewall.@defaults[0].forward='REJECT'
	set firewall.@defaults[0].forward='ACCEPT'
	set firewall.@defaults[0].synflood_protect='1'
	delete firewall.@zone[1].input='REJECT'
	delete firewall.@zone[1].forward='REJECT'
	set firewall.@zone[1].input='ACCEPT'
	set firewall.@zone[1].forward='ACCEPT'
	commit firewall
EOF

VDATE=$(opkg info tsmodem | grep Version | awk '{
    print substr($0,10,10)
}')

VNUM=$(opkg info tsmodem | grep Version | awk '{
    print substr($0,30,10)
}')

uci set tsmodem.version.number=$VNUM
uci set tsmodem.version.date=$VDATE
uci commit tsmodem

exit 0
