# luci/themes/luci-theme-qtech-1907/root/etc/uci-defaults/30_luci-theme-qtech
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Qtech=/luci-static/qtech
	set luci.main.mediaurlbase=/luci-static/qtech
	commit luci
EOF
exit 0
